# server/ava-beauty.service
[Unit]
Description=Ava Beauty Backend (Express + Prisma)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/ava-beauty
ExecStart=/usr/bin/node /var/www/ava-beauty/server/index.js
EnvironmentFile=/var/www/ava-beauty/.env
Environment=NODE_ENV=production
Environment=PORT=4000
Environment=JWT_SECRET=change_me_in_production
Environment=RL_WINDOW_MS=900000
Environment=RL_MAX=20
# Optimize Node.js memory usage (limit to 1GB to prevent OOM on small VPS)
Environment=NODE_OPTIONS="--max-old-space-size=1024"
# Optional initial superadmin seed vars (remove after first run)
# Environment=SUPERADMIN_USERNAME=admin
# Environment=SUPERADMIN_EMAIL=admin@local
# Environment=SUPERADMIN_PASSWORD=admin
Restart=on-failure
RestartSec=3s
# Ensure graceful stop
KillSignal=SIGINT
TimeoutStopSec=15
# Lower privileges if a service user exists (adjust as needed)
# User=www-data
# Group=www-data
# Journald logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
